# packaging/etc/systemd/system/abb-box64.service
[Unit]
Description=Synology Active Backup for Business Agent via Box64 on ARM64
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=root
Group=root
PIDFile=/run/synology-backupd.pid
ExecStartPre=/bin/sh -c 'modprobe synosnap'
ExecStart=/usr/local/bin/abb-box64-wrapper
ExecStop=/bin/kill -TERM $MAINPID
KillMode=control-group
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=30
Environment=BOX64_LOG=0

[Install]
WantedBy=multi-user.target
